# service/todoapi.service
[Unit]
Description=Gunicorn Daemon for FastAPI todo Application
After=network.target

[Service]
User=admin
Group=admin
WorkingDirectory=/home/admin/todoapi
ExecStart=/home/admin/todoapi/venv/bin/gunicorn -k uvicorn.workers.UvicornWorker -c /home/admin/todoapi/gunicorn_conf.py main:api

[Install]
WantedBy=multi-user.target
